# ops/systemd/uprava-node.service.example
[Unit]
Description=Uprava Node Daemon
After=network-online.target
Wants=network-online.target

[Service]
User=uprava
Group=uprava
EnvironmentFile=/etc/uprava/node.env
EnvironmentFile=-/opt/apps/uprava/.env.release
WorkingDirectory=/var/lib/uprava
ExecStart=/opt/apps/uprava/current/uprava-node
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/uprava /var/lib/uprava-node /srv/uprava-workspaces

[Install]
WantedBy=multi-user.target
